# systemd/claude-mcp.service
[Unit]
Description=Claude Bot MCP Server
After=network.target claude-discord.service
StartLimitBurst=5
StartLimitIntervalSec=60

[Service]
Type=simple
WorkingDirectory=/home/jason/projects/claude-bot
EnvironmentFile=/home/jason/projects/claude-bot/.env
Environment=NODE_ENV=production
Environment=PATH=/home/jason/.local/bin:/home/jason/projects/claude-bot/node_modules/.bin:/usr/local/bin:/usr/bin:/bin
ExecStart=/usr/bin/node /home/jason/projects/claude-bot/node_modules/tsx/dist/cli.mjs mcp/server.ts
KillMode=process
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=claude-mcp

[Install]
WantedBy=default.target
